# systemcheck.timer
[Unit]
Description=systemcheck timer

[Timer]
OnUnitActiveSec=14400s
OnBootSec=300s

[Install]
WantedBy=timers.target
