# systemcheck.timer
[Unit]
Description=systemcheck

[Timer]
OnUnitActiveSec=14400s
OnBootSec=90s

[Install]
WantedBy=timers.target
