# systemcheck.timer
[Unit]
Description=systemcheck timer

[Timer]
OnBootSec=300s
OnUnitActiveSec=14400s

[Install]
WantedBy=timers.target
