# wifistatus.timer
[Unit]
Description=wifi status check timer

[Timer]
OnBootSec=330s
OnUnitActiveSec=14400s

[Install]
WantedBy=timers.target
